# deploy/systemd/cloudflared-board.service
[Unit]
Description=Cloudflare Tunnel for Ops Dashboard (board.agenticker.cc -> 127.0.0.1:8000)
After=network.target
# 隧道断 = 看板公网不可达；崩溃发飞书告警（与 data-hub 同一 OnFailure 机制）。
OnFailure=data-onfailure@%n.service

[Service]
Type=simple
# cloudflared 与 uv 同装在 %h/.local/bin；systemd 无 PATH，必须绝对路径。
Environment="PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# --no-autoupdate 是顶层 flag，必须放在 tunnel 之前（否则报 flag not defined）。
ExecStart=%h/.local/bin/cloudflared --no-autoupdate tunnel --config %h/code/crossborder-ops-data-hub/deploy/cloudflared/config.yml run board
# SIGTERM(143) 是正常停止/重启退出码（同 data-hub.service），不加会被判 failed 误告警。
SuccessExitStatus=143 SIGTERM
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
