# deploy/systemd/data-refresh-tokens.timer
[Unit]
Description=жҜҸ 6 е°Ҹж—¶еҲ·ж–° TikTok Token
[Timer]
OnCalendar=*-*-* 00,06,12,18:41:00
Persistent=true
[Install]
WantedBy=timers.target
