# deploy/systemd/data-push-replenishment.timer
[Unit]
Description=每日推送补货采购单（上午 09:30，运营到岗即见）
[Timer]
OnCalendar=*-*-* 09:30:00
Persistent=true
[Install]
WantedBy=timers.target
