# deploy/systemd/data-sync-fulfillments.timer
[Unit]
Description=жҜҸ 15 еҲҶй’ҹи·‘ TikTok еҫ…еҸ‘иҙ§еҗҢжӯҘ
[Timer]
OnCalendar=*-*-* *:05,20,35,50:00
Persistent=true
[Install]
WantedBy=timers.target
