# deploy/systemd/data-sync-orders.timer
[Unit]
Description=Т»Ј 15 тѕєжњЪУиЉ TikTok У«бтЇЋтљїТГЦ
[Timer]
OnCalendar=*-*-* *:08,23,38,53:00
Persistent=true
[Install]
WantedBy=timers.target
