# deploy/systemd/data-sync-inventory.timer
[Unit]
Description=жҜҸ 15 еҲҶй’ҹи·‘ TikTok еә“еӯҳеҗҢжӯҘ
[Timer]
OnCalendar=*-*-* *:02,17,32,47:00
Persistent=true
[Install]
WantedBy=timers.target
